# elkabot.service
[Unit]
Description=MyElkaBot
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python /home/bots/elka_bot/main.py
Restart=always

[Install]
WantedBy=multi-user.target
